# install/supervisord.conf

[program:webhooks-tool]
command=###INSTALL_PATH###/index.js
autorestart=true
environment=NODE_ENV=production
user=webhooks
directory=###INSTALL_PATH###
